# crontab/reboot.sh
#!/bin/bash
echo "Khoi dong lai vps"
sudo reboot
